# _systemd-units/etc/systemd/user/polkit-agent.service
[Unit]
Description=Polkit Agent (polkit-dumb-agent)

[Service]
ExecStart=polkit-dumb-agent

[Install]
WantedBy=default.target
